# scripts/com.resellbot.telegram-poll-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.resellbot.telegram-poll-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/morganlynn/repos/resell-bot/.venv/bin/python3</string>
        <string>/Users/morganlynn/repos/resell-bot/scripts/telegram_poll_bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/morganlynn/repos/resell-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/morganlynn/repos/resell-bot/logs/poll-bot-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/morganlynn/repos/resell-bot/logs/poll-bot-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
